# src/service/innasrv.upstart

start on runlevel [2345]
stop on runlevel [016]

exec /usr/local/sbin/innasrv
